# deploy/upstart/app.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

script
  exec sudo -E -u username bash -c "/path/to/node /path/to/rejoice -c /path/to/hapi.cfg -p /path/to/modules >> /path/to/logfile 2&>1"
end script
